# mebtty.service
[Unit]
Description=MebTTY - Self-hosted web terminal
Documentation=https://github.com/mill413/mebtty
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/lib/mebtty
ExecStart=/usr/local/bin/mebtty
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Environment defaults (override in /etc/mebtty/mebtty.env)
Environment=MEBTTY_HOST=0.0.0.0
Environment=MEBTTY_PORT=18888
EnvironmentFile=-/etc/mebtty/mebtty.env

# Data directories
StateDirectory=mebtty
LogsDirectory=mebtty
Environment=MEBTTY_DATABASE_URL=sqlite+aiosqlite:////var/lib/mebtty/mebtty.db
Environment=MEBTTY_UPLOAD_DIR=/var/lib/mebtty/uploads
Environment=MEBTTY_PLUGIN_DIR=/var/lib/mebtty/plugins

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mebtty

# Security hardening
# MebTTY starts as root so it can authenticate local users and drop terminal
# sessions to the selected user's uid/gid, similar to sshd/login.
NoNewPrivileges=false
ProtectSystem=strict
# PAM modules such as pam_faillock need to update runtime auth state.
ReadWritePaths=/var/lib/mebtty /home /run/faillock /var/run/faillock
ProtectHome=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
