# pkg/deb/DEBIAN/postinst
#!/usr/bin/env bash
set -e

DATA_DIR="/var/lib/mebtty"
ENV_FILE="/etc/mebtty/mebtty.env"

append_commented_env() {
    local env_file="$1"
    local key="$2"
    local description="$3"
    local value="$4"

    if grep -Eq "^[[:space:]#]*${key}=" "$env_file"; then
        return
    fi

    {
        printf '\n# %s\n' "$description"
        printf '# %s=%s\n' "$key" "$value"
    } >> "$env_file"
}

append_optional_env() {
    local env_file="$1"

    append_commented_env "$env_file" \
        "MEBTTY_DATABASE_URL" \
        "Database connection string." \
        "sqlite+aiosqlite:////var/lib/mebtty/mebtty.db"
    append_commented_env "$env_file" \
        "MEBTTY_UPLOAD_DIR" \
        "Directory for uploaded files and avatars." \
        "/var/lib/mebtty/uploads"
    append_commented_env "$env_file" \
        "MEBTTY_BROWSE_ROOT" \
        "Root directory exposed by the file browser." \
        "/var/lib/mebtty"
    append_commented_env "$env_file" \
        "MEBTTY_STATIC_DIR" \
        "Path to frontend build output. Leave empty to auto-detect bundled assets." \
        ""
    append_commented_env "$env_file" \
        "MEBTTY_ACCESS_TOKEN_EXPIRE_MINUTES" \
        "JWT access token lifetime in minutes." \
        "60"
    append_commented_env "$env_file" \
        "MEBTTY_REFRESH_TOKEN_EXPIRE_DAYS" \
        "JWT refresh token lifetime in days." \
        "7"
    append_commented_env "$env_file" \
        "MEBTTY_MAX_UPLOAD_SIZE" \
        "Maximum upload size in bytes." \
        "104857600"
    append_commented_env "$env_file" \
        "MEBTTY_REGISTRATION_ENABLED" \
        "Allow creating new MebTTY web accounts." \
        "true"
    append_commented_env "$env_file" \
        "MEBTTY_ALLOW_ROOT_LOCAL_USER" \
        "Allow selecting root as a local terminal user." \
        "false"
    append_commented_env "$env_file" \
        "MEBTTY_ALGORITHM" \
        "JWT signing algorithm." \
        "HS256"
    append_commented_env "$env_file" \
        "MEBTTY_REDIS_URL" \
        "Reserved for deployments that enable Redis-backed features." \
        "redis://localhost:6379/0"
    append_commented_env "$env_file" \
        "MEBTTY_DEV_FRONTEND_HOST" \
        "Local deploy.sh development only; ignored by packaged systemd service." \
        "127.0.0.1"
    append_commented_env "$env_file" \
        "MEBTTY_DEV_FRONTEND_PORT" \
        "Local deploy.sh development only; ignored by packaged systemd service." \
        "3000"
}

# Create data directories
mkdir -p "$DATA_DIR/uploads"
chown -R root:root "$DATA_DIR"
chmod 750 "$DATA_DIR" "$DATA_DIR/uploads"

# Generate environment config with random secret key if not exists
if [ ! -f "$ENV_FILE" ]; then
    install -d -m 750 -o root -g root "$(dirname "$ENV_FILE")"
    SECRET=$(python3 -c 'import secrets; print(secrets.token_hex(32))' 2>/dev/null || openssl rand -hex 32)
    cat > "$ENV_FILE" <<EOF
# MebTTY environment configuration
# Modify and restart: systemctl restart mebtty

# Required: JWT signing key. Keep this secret and stable across restarts.
MEBTTY_SECRET_KEY=$SECRET

# Server bind address.
MEBTTY_HOST=0.0.0.0

# Server listen port.
MEBTTY_PORT=18888

# PAM service used to verify local user passwords.
MEBTTY_PAM_SERVICE=login

# Database connection string.
# MEBTTY_DATABASE_URL=sqlite+aiosqlite:////var/lib/mebtty/mebtty.db

# Directory for uploaded files and avatars.
# MEBTTY_UPLOAD_DIR=/var/lib/mebtty/uploads

# Root directory exposed by the file browser.
# MEBTTY_BROWSE_ROOT=/var/lib/mebtty

# Path to frontend build output. Leave empty to auto-detect bundled assets.
# MEBTTY_STATIC_DIR=

# JWT access token lifetime in minutes.
# MEBTTY_ACCESS_TOKEN_EXPIRE_MINUTES=60

# JWT refresh token lifetime in days.
# MEBTTY_REFRESH_TOKEN_EXPIRE_DAYS=7

# Maximum upload size in bytes.
# MEBTTY_MAX_UPLOAD_SIZE=104857600

# Allow creating new MebTTY web accounts.
# MEBTTY_REGISTRATION_ENABLED=true

# Allow selecting root as a local terminal user.
# MEBTTY_ALLOW_ROOT_LOCAL_USER=false

# JWT signing algorithm.
# MEBTTY_ALGORITHM=HS256

# Reserved for deployments that enable Redis-backed features.
# MEBTTY_REDIS_URL=redis://localhost:6379/0

# Local deploy.sh development only; ignored by packaged systemd service.
# MEBTTY_DEV_FRONTEND_HOST=127.0.0.1
# MEBTTY_DEV_FRONTEND_PORT=3000
EOF
    chown root:root "$ENV_FILE"
    chmod 640 "$ENV_FILE"
else
    chown root:root "$(dirname "$ENV_FILE")" "$ENV_FILE"
    chmod 750 "$(dirname "$ENV_FILE")"
    append_optional_env "$ENV_FILE"
    chmod 640 "$ENV_FILE"
fi

# Enable and start systemd service
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload
    systemctl enable mebtty.service
    systemctl start mebtty.service
fi

echo ""
echo "  MebTTY installed successfully!"
echo ""
echo "  Config:     $ENV_FILE"
echo "  Data:       $DATA_DIR"
echo "  Service:    root user"
echo "  Web UI:     http://localhost:18888"
echo ""
echo "  systemctl start mebtty"
echo "  systemctl stop mebtty"
echo "  systemctl status mebtty"
echo "  journalctl -u mebtty -f"
echo ""
