# pkg/deb/DEBIAN/postinst
#!/usr/bin/env bash
set -e

DATA_DIR="/var/lib/mebtty"
ENV_FILE="/etc/mebtty/mebtty.env"

# Create data directories
mkdir -p "$DATA_DIR/uploads"
chown -R root:root "$DATA_DIR"
chmod 750 "$DATA_DIR" "$DATA_DIR/uploads"

# Generate environment config with random secret key if not exists
if [ ! -f "$ENV_FILE" ]; then
    install -d -m 750 -o root -g root "$(dirname "$ENV_FILE")"
    SECRET=$(python3 -c 'import secrets; print(secrets.token_hex(32))' 2>/dev/null || openssl rand -hex 32)
    cat > "$ENV_FILE" <<EOF
# MebTTY environment configuration
# Modify and restart: systemctl restart mebtty

MEBTTY_SECRET_KEY=$SECRET
MEBTTY_HOST=0.0.0.0
MEBTTY_PORT=18888
MEBTTY_PAM_SERVICE=login
EOF
    chown root:root "$ENV_FILE"
    chmod 640 "$ENV_FILE"
else
    chown root:root "$(dirname "$ENV_FILE")" "$ENV_FILE"
    chmod 750 "$(dirname "$ENV_FILE")"
    chmod 640 "$ENV_FILE"
fi

# Enable and start systemd service
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload
    systemctl enable mebtty.service
    systemctl start mebtty.service
fi

echo ""
echo "  MebTTY installed successfully!"
echo ""
echo "  Config:     $ENV_FILE"
echo "  Data:       $DATA_DIR"
echo "  Service:    root user"
echo "  Web UI:     http://localhost:18888"
echo ""
echo "  systemctl start mebtty"
echo "  systemctl stop mebtty"
echo "  systemctl status mebtty"
echo "  journalctl -u mebtty -f"
echo ""
